# systemd/xplus-contacts.service
[Unit]
Description=XPlus Contacts — Personal Contact Manager
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=%INSTALL_USER%
WorkingDirectory=%PROJECT_DIR%
ExecStart=%PROJECT_DIR%/.venv/bin/python -m src.apps.contacts.server
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
