# launchd/com.personal-ai-assistant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.personal-ai-assistant</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/ai-assistant/.venv/bin/python</string>
        <string>-m</string>
        <string>src.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/opt/ai-assistant</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/opt/ai-assistant/logs/launchd-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/opt/ai-assistant/logs/launchd-stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
